# faceauth.service
[Unit]
Description=FaceAuth Authentication Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/faceauthd
Environment=RUST_LOG=info
Restart=always
RestartSec=5
User=root
Group=root
# Prevent swapping — ensures fast wake-up after sleep
MemorySwapMax=0
# Run below normal scheduling priority so ONNX inference doesn't
# starve the compositor/desktop on resume from sleep.
Nice=10
CPUSchedulingPolicy=batch
IOSchedulingClass=best-effort
IOSchedulingPriority=6
# Reduce memory fragmentation
Environment=MALLOC_ARENA_MAX=2
# Cap ALL ONNX / OpenMP / MKL thread pools to 1.
# Without this, int8 ONNX models still spawn a thread per CPU core
# via OMP/MKL *even when* intra_threads=1 is set in the session,
# causing the CPU spike on every login and especially after sleep.
Environment=OMP_NUM_THREADS=1
Environment=MKL_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=ORT_DISABLE_SPINNING=1


[Install]
WantedBy=multi-user.target
