# faceauth.service
[Unit]
Description=FaceAuth Authentication Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/faceauthd
Environment=RUST_LOG=info
Restart=always
RestartSec=5
User=root
Group=root
# Performance Tuning
# Prevent swapping to ensure fast wake-up after sleep
MemorySwapMax=0
# Give slight CPU priority for authentication speed
Nice=-10
# Reduce memory fragmentation and virtual memory overhead
Environment=MALLOC_ARENA_MAX=2


[Install]
WantedBy=multi-user.target
